# backup.service.in
[Unit]
Description=Backup service

[Service]
Type=simple
ExecStart={{ libexedir }}/backup
