# backup.timer
[Unit]
Description=Backup service timer

[Timer]
OnCalendar=*-*-* 3:00:00
Persistent=yes

[Install]
WantedBy=timers.target
